# packaging/linux/systemd/vpn-server.service
[Unit]
# VPN 服务端 system 级 systemd service（区别于客户端 vpn-cli 的 user service）。
Description=VPN Server (vpn-server)
Documentation=https://github.com/Shangguanjunjie/vpn
# 数据平面需要网络就绪后再启动。
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vpn-server
Group=vpn-server
# 服务端配置全部走环境变量，由安装包写入的 EnvironmentFile 提供默认值。
EnvironmentFile=/etc/vpn-server/vpn-server.env
ExecStart=/usr/bin/vpn-server
Restart=on-failure
RestartSec=5

# 异地组网/站点互通：kernel 后端启动后开启 IPv4 转发并放行 wg0 互转发，
# 停止后清理 iptables 规则（脚本读取本单元注入的 env，幂等 + best-effort）。
# noop 后端下脚本会自行跳过。两个 ExecStartPost/ExecStopPost 需 root 能力，
# 这里以 root 权限运行（见下方 PermissionsStartOnly / 单元未降权前缀），
# 子进程不继承 User= 降权由 "+" 前缀保证。
ExecStartPost=+/usr/lib/vpn-server/wg-forward-setup.sh
ExecStopPost=+/usr/lib/vpn-server/wg-forward-teardown.sh

# 创建 WireGuard 接口 / 配置路由 / 改 sysctl 与 iptables 需要 CAP_NET_ADMIN。
# 注意：sysctl net.ipv4.ip_forward 写入需要 CAP_NET_ADMIN（且不能被 ProtectKernelTunables 阻断）。
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

# 数据目录（SQLite 数据库、JWT 密钥等）由 systemd 托管，归 vpn-server 用户所有。
StateDirectory=vpn-server
WorkingDirectory=/var/lib/vpn-server

# 基础加固（可按真机环境调整；TUN 设备访问保留）。
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
DeviceAllow=/dev/net/tun rw

[Install]
WantedBy=multi-user.target
